# infra/ec2/systemd/botboundary-model.service
[Unit]
Description=BotBoundary Model Service (FastAPI)
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/opt/botboundary/Model/login_auth
EnvironmentFile=/opt/botboundary/.env
# Bind to loopback so it is NOT publicly reachable.
ExecStart=/opt/botboundary/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
